# weiss/crontab/user-jobs
*/20 * * * * /home/weiss/rust/recentf/recentf clean
* * * * * env DISPLAY=:0 bb /home/weiss/clojure/scripts/xorg-set-value.clj
*/3 * * * * /home/weiss/rust/webman/webman-cli/provide.sh


